# x-ui-pro-install.service
[Unit]
Description=First-boot installer for x-ui-pro
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/etc/x-ui-pro/.installed
StartLimitIntervalSec=1200
StartLimitBurst=10

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/x-ui-pro.env
ExecStart=/usr/local/bin/install-wrapper.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=1800
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
